# samples/avorf1com.supervisor.conf
[program:avorf1com]
command=/opt/avorf1com/current/bin/avorf1com
user=root
autostart=true
autorestart=true